# systemd/start-smb-containers.service
[Unit]
Description=Start SMB/NFS Containers after TrueNAS
Documentation=https://github.com/yourusername/proxmox-smb-autostart
After=network.target network-online.target pve-storage.target pve-guests.target
Requires=network-online.target
Wants=pve-storage.target pve-guests.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/start-smb-containers.sh
TimeoutStartSec=400
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smb-autostart

[Install]
WantedBy=multi-user.target
